# infra/systemd-services/recipeyak-algolia_sync_continuous.service
[Unit]
Description=Runs algolia sync job once
[Service]
Restart=always
ExecStart=/usr/bin/docker run \
        --rm \
        --network host \
        --log-driver=journald \
        --env-file=/root/.env-production \
        recipeyak/django:{{GIT_SHA}} \
        ./.venv/bin/python -m recipeyak.jobs.live_search_sync --no-single-run
KillSignal=SIGINT
[Install]
WantedBy=multi-user.target
